# services/ipsum/crontab
# update blocklist every day
0 2 * * * /app/ipsum.sh